# crontab
# Запуск скрипта раз в сутки в 00:00
0 0 * * * cd /app && /usr/local/bin/python /app/fetch_coordinates.py >> /var/log/cron.log 2>&1

#212
# Пустая строка в конце файла обязательна 

